# dot_resources/etc/supervisor/services.d/nginx.ini
[program:nginx]
; environment = LD_LIBRARY_PATH="/usr/local/openresty/openssl111/lib:$LD_LIBRARY_PATH"
# command=%(ENV_HOME)s/.local/openresty/nginx/sbin/nginx -c %(ENV_HOME)s/.local/etc/nginx/nginx.conf -p %(ENV_HOME)s/.local/openresty/nginx/ -g "daemon off;"
# command=/usr/local/openresty/nginx/sbin/nginx -c %(ENV_HOME)s/.local/etc/nginx/nginx.conf -p /usr/local/openresty/nginx/ -g "daemon off;"
command=/usr/local/openresty/nginx/sbin/nginx -c %(ENV_HOME)s/.local/etc/nginx/nginx.conf -g "daemon off;"


directory=%(ENV_HOME)s/.local/openresty/lualib
priority=50
redirect_stderr=false
;stdout_logfile=/tmp/supervisor/%(program_name)s.log ; log logs into file
# stdout_logfile=/proc/1/fd/1 ; log error into stdout to see in docker logs
stdout_logfile_maxbytes=1MB
stdout_logfile=%(ENV_LOG_PATH)s/supervisor/%(program_name)s-stdout-%(process_num)s-%(ENV_USER)s-%(host_node_name)s.log ; log logs into file
stdout_capture_maxbytes=1MB
stderr_logfile_maxbytes=1MB
stderr_logfile=%(ENV_LOG_PATH)s/supervisor/%(program_name)s-stderr-%(process_num)s-%(ENV_USER)s-%(host_node_name)s.log ; log logs into file
stderr_capture_maxbytes=1MB
autostart=true ; start at supervisord start (default: true)
autorestart=true ; whether/when to restart (default: unexpected) 
# stdout_logfile=/tmp/supervisor/%(program_name)s.log ; log logs into file
# stdout_logfile=./%(program_name)s_%(ENV_HOST)s.log ; log logs into file
startretries=5   ; max # of serial start failures (default 3)
# environment=PATH="%(ENV_HOME)s/.conda/bin:%(ENV_PATH)s",VNC_PW="automated"
# environment=PATH="%(ENV_HOME)s/.conda/bin:%(ENV_PATH)s",LD_LIBRARY_PATH=%(ENV_HOME)s/.local/openresty/openssl111/lib:%(ENV_HOME)s/.local/openresty/pcre/lib/:%(ENV_HOME)s/.local/openresty/luajit/lib${LD_LIBRARY_PATH:+:$LD_LIBRARY_PATH},PATH="%(ENV_HOME)s/.local/openresty/bin:%(ENV_HOME)s/.local/openresty/nginx/sbin:$PATH"
# environment=PATH="%(ENV_HOME)s/.conda/bin:%(ENV_PATH)s",LD_LIBRARY_PATH="%(ENV_HOME)s/.local/openresty/openssl111/lib:%(ENV_HOME)s/.local/openresty/pcre/lib/:%(ENV_HOME)s/.local/openresty/luajit/lib${LD_LIBRARY_PATH:+:%(ENV_LD_LIBRARY_PATH)s}",PATH="%(ENV_HOME)s/.local/openresty/bin:%(ENV_HOME)s/.local/openresty/nginx/sbin:%(ENV_PATH)s"
# environment=LD_LIBRARY_PATH="%(ENV_HOME)s/.local/openresty/openssl111/lib:%(ENV_HOME)s/.local/openresty/pcre/lib/:%(ENV_HOME)s/.local/openresty/luajit/lib:%(ENV_LD_LIBRARY_PATH)s",PATH="%(ENV_HOME)s/.local/openresty/bin:%(ENV_HOME)s/.local/openresty/nginx/sbin:%(ENV_PATH)s"
# ENV LUA_PATH="/usr/local/openresty/site/lualib/?.ljbc;/usr/local/openresty/site/lualib/?/init.ljbc;/usr/local/openresty/lualib/?.ljbc;/usr/local/openresty/lualib/?/init.ljbc;/usr/local/openresty/site/lualib/?.lua;/usr/local/openresty/site/lualib/?/init.lua;/usr/local/openresty/lualib/?.lua;/usr/local/openresty/lualib/?/init.lua;./?.lua;/usr/local/openresty/luajit/share/luajit-2.1.0-beta3/?.lua;/usr/local/share/lua/5.1/?.lua;/usr/local/share/lua/5.1/?/init.lua;/usr/local/openresty/luajit/share/lua/5.1/?.lua;/usr/local/openresty/luajit/share/lua/5.1/?/init.lua"

# ENV LUA_CPATH="/usr/local/openresty/site/lualib/?.so;/usr/local/openresty/lualib/?.so;./?.so;/usr/local/lib/lua/5.1/?.so;/usr/local/openresty/luajit/lib/lua/5.1/?.so;/usr/local/lib/lua/5.1/loadall.so;/usr/local/openresty/luajit/lib/lua/5.1/?.so"

# environment=PATH=%(ENV_HOME)s/.local/openresty/luajit/bin:%(ENV_HOME)s/.local/openresty/nginx/sbin:%(ENV_HOME)s/.local/openresty/bin:%(ENV_PATH)s,RESTY_PREFIX=%(ENV_HOME)s/.local/openresty,NGINX_PREFIX=%(ENV_HOME)s/.local/openresty/nginx,LUA_PATH="%(ENV_HOME)s/.local/openresty/site/lualib/?.ljbc;%(ENV_HOME)s/.local/openresty/site/lualib/?/init.ljbc;%(ENV_HOME)s/.local/openresty/lualib/?.ljbc;%(ENV_HOME)s/.local/openresty/lualib/?/init.ljbc;%(ENV_HOME)s/.local/openresty/site/lualib/?.lua;%(ENV_HOME)s/.local/openresty/site/lualib/?/init.lua;%(ENV_HOME)s/.local/openresty/lualib/?.lua;%(ENV_HOME)s/.local/openresty/lualib/?/init.lua;./?.lua;%(ENV_HOME)s/.local/openresty/luajit/share/luajit-2.1.0-beta3/?.lua;/usr/local/share/lua/5.1/?.lua;/usr/local/share/lua/5.1/?/init.lua;%(ENV_HOME)s/.local/openresty/luajit/share/lua/5.1/?.lua;%(ENV_HOME)s/.local/openresty/luajit/share/lua/5.1/?/init.lua",LUA_CPATH="%(ENV_HOME)s/.local/openresty/site/lualib/?.so;%(ENV_HOME)s/.local/openresty/lualib/?.so;./?.so;/usr/local/lib/lua/5.1/?.so;%(ENV_HOME)s/.local/openresty/luajit/lib/lua/5.1/?.so;/usr/local/lib/lua/5.1/loadall.so;%(ENV_HOME)s/.local/openresty/luajit/lib/lua/5.1/?.so",LD_LIBRARY_PATH="%(ENV_HOME)s/.local/openresty/openssl111/lib:%(ENV_HOME)s/.local/openresty/pcre/lib/:%(ENV_HOME)s/.local/openresty/luajit/lib"
environment=PATH="/usr/local/openresty/luajit/bin:/usr/local/openresty/nginx/sbin:/usr/local/openresty/bin:%(ENV_PATH)s",RESTY_PREFIX="/usr/local/openresty",NGINX_PREFIX="/usr/local/openresty/nginx",LUA_PATH="/usr/local/openresty/site/lualib/?.ljbc;/usr/local/openresty/site/lualib/?/init.ljbc;/usr/local/openresty/lualib/?.ljbc;/usr/local/openresty/lualib/?/init.ljbc;/usr/local/openresty/site/lualib/?.lua;/usr/local/openresty/site/lualib/?/init.lua;/usr/local/openresty/lualib/?.lua;/usr/local/openresty/lualib/?/init.lua;./?.lua;/usr/local/openresty/luajit/share/luajit-2.1.0-beta3/?.lua;/usr/local/share/lua/5.1/?.lua;/usr/local/share/lua/5.1/?/init.lua;/usr/local/openresty/luajit/share/lua/5.1/?.lua;/usr/local/openresty/luajit/share/lua/5.1/?/init.lua",LUA_CPATH="/usr/local/openresty/site/lualib/?.so;/usr/local/openresty/lualib/?.so;./?.so;/usr/local/lib/lua/5.1/?.so;/usr/local/openresty/luajit/lib/lua/5.1/?.so;/usr/local/lib/lua/5.1/loadall.so;/usr/local/openresty/luajit/lib/lua/5.1/?.so",LD_LIBRARY_PATH="/usr/local/openresty/openssl111/lib:/usr/local/openresty/pcre/lib/:/usr/local/openresty/luajit/lib"



# export PATH=/usr/local/openresty/luajit/bin:/usr/local/openresty/nginx/sbin:/usr/local/openresty/bin:${PATH}
# export RESTY_PREFIX=/usr/local/openresty
# export NGINX_PREFIX=/usr/local/openresty/nginx
# export LUA_PATH="/usr/local/openresty/site/lualib/?.ljbc;/usr/local/openresty/site/lualib/?/init.ljbc;/usr/local/openresty/lualib/?.ljbc;/usr/local/openresty/lualib/?/init.ljbc;/usr/local/openresty/site/lualib/?.lua;/usr/local/openresty/site/lualib/?/init.lua;/usr/local/openresty/lualib/?.lua;/usr/local/openresty/lualib/?/init.lua;./?.lua;/usr/local/openresty/luajit/share/luajit-2.1.0-beta3/?.lua;/usr/local/share/lua/5.1/?.lua;/usr/local/share/lua/5.1/?/init.lua;/usr/local/openresty/luajit/share/lua/5.1/?.lua;/usr/local/openresty/luajit/share/lua/5.1/?/init.lua"
# export LUA_CPATH="/usr/local/openresty/site/lualib/?.so;/usr/local/openresty/lualib/?.so;./?.so;/usr/local/lib/lua/5.1/?.so;/usr/local/openresty/luajit/lib/lua/5.1/?.so;/usr/local/lib/lua/5.1/loadall.so;/usr/local/openresty/luajit/lib/lua/5.1/?.so"
# export LD_LIBRARY_PATH="/usr/local/openresty/openssl111/lib:/usr/local/openresty/pcre/lib/:/usr/local/openresty/luajit/lib"

