# private_dot_local/etc/supervisor/services.d/code.ini
[program:code]
command=code tunnel
autostart=true
autorestart=true